# deploy/systemd/coze-js-api-container-cleanup.timer
[Unit]
Description=Run coze-js-api stopped-container cleanup hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=coze-js-api-container-cleanup.service

[Install]
WantedBy=timers.target
